# linux-os/xrdp/supervisor.d/xrdp.conf
[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
user=root
autorestart=true
priority=400
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
